# myimpression.service
[Unit]
Description=MyImpression Display Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=gavin
Group=gavin
WorkingDirectory=/home/gavin/MyImpression
Environment=PATH=/home/gavin/MyImpression/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/gavin/MyImpression/venv/bin/python /home/gavin/MyImpression/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# GPIO access
SupplementaryGroups=gpio

[Install]
WantedBy=multi-user.target
